# openclaw-village.service
[Unit]
Description=OpenClaw Village Orchestrator
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/root/openclaw-cloud/village
EnvironmentFile=/root/openclaw-cloud/village/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-village

# Graceful shutdown: send SIGTERM, wait 30s for tick to finish
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
